# board/raspi/pi0/rootfs-overlay/etc/systemd/system/telnetd.service
[Service]
ExecStart=/usr/sbin/telnetd -F

[Install]
WantedBy=network.target
